# example/supervisor.d/consul.ini
[program:consul]
command=/usr/local/bin/consul-entrypoint.sh consul agent -server -bootstrap-expect=1 -data-dir=/consul/data
autostart=true
startsecs=0
autorestart=true
stderr_logfile=/var/log/consul_stderr.log
stdout_logfile=/var/log/consul_stdout.log