# example/supervisor.d/insert.ini
[program:testdata]
command=/itko/testdata/insert.sh
autostart=false
startsecs=45
autorestart=false
stderr_logfile=/var/log/add-leaves_stderr.log
stdout_logfile=/var/log/add-leaves_stdout.log