# example/supervisor.d/insert.ini
[program:testdata]
command=/itko/testdata/insert.sh
autostart=true
startsecs=0
autorestart=false
stderr_logfile=/var/log/testdata_insert_stderr.log
stdout_logfile=/var/log/testdata_insert_stdout.log