# example/supervisor.d/itko-submit.ini
[program:itko-submit]
command=/bin/sh -l -c "sleep 8 && /itko/start-itko-submit.sh"
autostart=true
autorestart=true
stderr_logfile=/var/log/itko-submit_stderr.log
stdout_logfile=/var/log/itko-submit_stdout.log