# example/supervisor.d/itko-monitor.ini
[program:itko-monitor]
command=/itko/itko-monitor -listen-address ${ITKO_MONITOR_LISTEN_ADDRESS}:${ITKO_MONITOR_LISTEN_PORT} -mask-size ${ITKO_MASK_SIZE} -store-directory ${ITKO_ROOT_DIRECTORY}
autostart=true
autorestart=true
stderr_logfile=/var/log/itko-monitor_stderr.log
stdout_logfile=/var/log/itko-monitor_stdout.log