# example/supervisor.d/itko-submit.ini
[program:itko-submit]
command=/itko/itko-submit -getentries_metrics -kv-path "${ITKO_KV_PATH}" -listen-address "${ITKO_SUBMIT_LISTEN_ADDRESS}:${ITKO_SUBMIT_LISTEN_PORT}"
autostart=true
autorestart=true
stderr_logfile=/var/log/itko-submit_stderr.log
stdout_logfile=/var/log/itko-submit_stdout.log