# example/supervisor.d/itko-submit.ini
[program:itko-submit]
command=/itko/start-itko-submit.sh
autostart=true
startsecs=30
autorestart=true
stderr_logfile=/var/log/itko-submit_stderr.log
stdout_logfile=/var/log/itko-submit_stdout.log