# example/supervisor.d/testdata.ini
[program:testdata]
command=/itko/testdata/testdata.sh
autostart=true
startsecs=0
autorestart=false
stderr_logfile=/var/log/testdata_stderr.log
stdout_logfile=/var/log/testdata_stdout.log